# dotconfig/systemd/user/display-hz-watch.service
[Unit]
Description=Adjust display Hz on charge state change
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.config/hypr/display-hz-watch.sh
Restart=on-failure

[Install]
WantedBy=graphical-session.target
